# case/backlight.service
[Unit]
Description=Backlight Daemon

[Service]
ExecStart=/opt/backlight.sh
Restart=always
RestartSec=3
User=pi

[Install]
WantedBy=default.target
